# .docker/services.d/esb/run
#!/usr/bin/with-contenv sh
cd /service
sudo -E -u www-data bin/console esb:pool start -f